# pktables.timer
[Unit]
Description=Run pktables.service every once in a while

[Timer]
OnCalendar=*-*-* *:00:00
